# urbana.cron
0 * * * * /usr/local/bin/python3 /Users/david/projects/urbana/urbana-scout.py
